# centos_marytts/supervisord.conf
[supervisord]
nodaemon = true

[program:marytts]
command = /marytts/marytts server 
logfile = /var/log/supervisor/maryttsserver.log
autorestart=true